# supervisord.conf
[supervisord]
nodaemon=true

[program:http-server]
command=http-server ./frontend/build

[program:python3]
command=python3 ./backend/python-api/wsgi.py